# com.youxi.puying.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.youxi.puying.server</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/chuningwang/.nvm/versions/node/v22.12.0/bin/node</string>
        <string>/Users/chuningwang/WorkBuddy/20260504230058/server.cjs</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/chuningwang/WorkBuddy/20260504230058</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_OPTIONS</key>
        <string></string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/chuningwang/WorkBuddy/20260504230058/logs/server-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/chuningwang/WorkBuddy/20260504230058/logs/server-launchd.err</string>
    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
